# gpt/XML/61672.xml
<service_bundle type='manifest' name='myapp'>\n       <service name='site/myapp' type='service'>\n           <create_default_instance enabled='true' />\n           <single_instance />\n           <dependency name='network' grouping='require_all' restart_on='none' type='service'>\n               <service_fmri value='svc:/network/physical:default' />\n           </dependency>\n           <exec_method type='method' name='start' exec='/path/to/start/script' timeout_seconds='60' />\n           <exec_method type='method' name='stop' exec='/path/to/stop/script' timeout_seconds='60' />\n           <property_group name='startd' type='framework'>\n               <propval name='duration' type='astring' value='transient' />\n           </property_group>\n           <stability value='Evolving' />\n       </service>\n   </service_bundle>