# llama31/INI/25681.ini
[program:my-php-service]\nprocess_name=my-php-service\ncommand=/usr/bin/php /path/to/your/script.php\nautostart=true\nautorestart=true