# llama31/INI/25738.ini
[program:worker]\nprocess_name=worker\ncommand=php /path/to/worker/script.php\nautostart=true\nautorestart=true