# gpt/INI/25406.ini
[program:php-worker]\n     command=/usr/bin/php /path/to/your/worker.php\n     autostart=true\n     autorestart=true\n     user=youruser\n     numprocs=5\n     stdout_logfile=/var/log/supervisor/php-worker.log\n     stderr_logfile=/var/log/supervisor/php-worker_error.log