# gpt/INI/25719.ini
[program:php-daemon]\n   command=php /path/to/your/script.php\n   autostart=true\n   autorestart=true\n   stderr_logfile=/var/log/php-daemon.err.log\n   stdout_logfile=/var/log/php-daemon.out.log