# qwen25-coder/INI/20665.ini
[program:php_worker]\n   command=/usr/bin/php /path/to/your/worker.php\n   autostart=true\n   autorestart=true\n   stderr_logfile=/var/log/php_worker.err.log\n   stdout_logfile=/var/log/php_worker.out.log\n   numprocs=5  ; Number of worker processes to ru