# qwen25-coder/INI/20680.ini
[program:your_script]\n    command=/usr/bin/php /path/to/your/script.php\n    autostart=true\n    autorestart=true\n    stderr_logfile=/var/log/your_script.err.log\n    stdout_logfile=/var/log/your_script.out.log